# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
 
[program:sshd]
command=/root/ssh_init.sh
 
[program:vncserver]
command=/root/vnc_init.sh
